# ciserver/cookbooks/jenkins/templates/default/jenkins.conf.erb
description "Jenkins CI server"

start on runlevel [2345]
stop on runlevel [!2345]
#Send KILL after 5 seconds
kill timeout 5
respawn

env JENKINS_LIB_DIR="/var/lib/jenkins"
env JENKINS_LOG_DIR="/var/log/jenkins"
env JENKINS_RUN_DIR="/var/run/jenkins"
env WEB_PORT=<%= @webport %>
env MAVEN_HOME=/usr/local/maven
env ANT_HOME=/usr/share/ant
export MAVEN_HOME
export ANT_HOME

setuid jenkins

script
exec java -jar $JENKINS_LIB_DIR/jenkins.war --httpPort=$WEB_PORT --webroot=$JENKINS_RUN_DIR --ajp13Port=-1 --logfile=$JENKINS_LOG_DIR/jenkins.log
end script

